# supervisor.conf
; =============================================================================
; RouterVault — Supervisor Queue Worker Configuration
; =============================================================================
; Install on production:
;   1. Copy to /etc/supervisor/conf.d/routervault.conf
;   2. sudo supervisorctl reread
;   3. sudo supervisorctl update
;   4. sudo supervisorctl start routervault-queue:*
;
; Restart after deploy:
;   sudo supervisorctl restart routervault-queue:*
;
; View status:
;   sudo supervisorctl status routervault-queue:*
; =============================================================================

[program:routervault-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/RouterVault/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-jobs=100 --max-time=3600 --timeout=240 --rest=1
numprocs=4
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=300
stopasgroup=true
killasgroup=true
user=www-data
directory=/var/www/html/RouterVault
stdout_logfile=/var/www/html/RouterVault/storage/logs/supervisor-queue-out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/www/html/RouterVault/storage/logs/supervisor-queue-err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
redirect_stderr=true
